# rootfs/etc/init/tf700.conf
# tf700

description  "Asus tf700t specific init script"

start on started udev

task

script
    # Create ssh host keys
    if [ ! -e /etc/ssh/ssh_host_dsa_key \
        -o -e /etc/ssh/ssh_host_rsa_key \
        -o -e /etc/ssh/ssh_host_ecdsa_key ]; then
            dpkg-reconfigure openssh-server
    fi

    # MMC read ahead size
    if [ -e /sys/block/mmcblk0/queue/read_ahead_kb ]; then
        echo 2048 > /sys/block/mmcblk0/queue/read_ahead_kb
    fi
    if [ -e /sys/block/mmcblk1/queue/read_ahead_kb ]; then
        echo 2048 > /sys/block/mmcblk1/queue/read_ahead_kb
    fi

    # Disable android lowmemkiller
    if [ -e /sys/module/lowmemorykiller/parameters/minfree ]; then
        echo "0,0,0,0,0,0" > /sys/module/lowmemorykiller/parameters/minfree
    fi

    # --- CPU ---
    # Block lowpower core
    #  it should be enabled on screen off or in strict powersawe mode
    [ -e /sys/module/cpu_tegra3/parameters/no_lp ] && echo 1 > /sys/module/cpu_tegra3/parameters/no_lp

    # Change conservative mode
    if [ -e /sys/devices/system/cpu/cpu0/cpufreq/scaling_governor ]; then
        echo conservative > /sys/devices/system/cpu/cpu0/cpufreq/scaling_governor
        [ -e /sys/devices/system/cpu/cpufreq/conservative/up_threshold ] && echo 95 > /sys/devices/system/cpu/cpufreq/conservative/up_threshold
        [ -e /sys/devices/system/cpu/cpufreq/conservative/down_threshold ] && echo 75 > /sys/devices/system/cpu/cpufreq/conservative/down_threshold
        [ -e /sys/devices/system/cpu/cpufreq/conservative/freq_step ] && echo 3 > /sys/devices/system/cpu/cpufreq/conservative/freq_step

        echo interactive > /sys/devices/system/cpu/cpu0/cpufreq/scaling_governor
        [ -e /sys/devices/system/cpu/cpufreq/interactive/boost_factor ] && echo 2 > /sys/devices/system/cpu/cpufreq/interactive/boost_factor
        [ -e /sys/devices/system/cpu/cpufreq/interactive/input_boost ] && echo 1 > /sys/devices/system/cpu/cpufreq/interactive/input_boost
        [ -e /sys/devices/system/cpu/cpufreq/interactive/sustain_load ] && echo 80 > /sys/devices/system/cpu/cpufreq/interactive/sustain_load
    fi

    # CPU hotplugging
    [ -e /sys/module/cpu_tegra3/parameters/auto_hotplug ] && echo 1 > /sys/module/cpu_tegra3/parameters/auto_hotplug

    # lp2 idle state
    [ -e /sys/module/cpuidle/parameters/lp2_in_idle ] && echo "Y" > /sys/module/cpuidle/parameters/lp2_in_idle

    # --- GPS ---
    if [ ! -e /sys/block/mmcblk0/queue/read_ahead_kb ]; then
        [ -e /sys/class/gpio/export ] && echo 162 > /sys/class/gpio/export
        [ -e /sys/class/gpio/gpio162/value ] && echo 0 > /sys/class/gpio/gpio162/value
        [ -e /sys/class/gpio/gpio162/direction ] && echo out > /sys/class/gpio/gpio162/direction
    fi

    # --- USB ---
    # Wakeup on usb insert
    [ -e /sys/bus/usb/devices/1-1/power/wakeup ] && echo enabled > /sys/bus/usb/devices/1-1/power/wakeup


    # --- Mount android ---
    [ -d /mnt/android/system ] || mkdir -p /mnt/android/system
    [ -d /mnt/android/data ]   || mkdir -p /mnt/android/data
    # Mount system
    mount | grep '/dev/mmcblk0p1' || mount /dev/mmcblk0p1 /mnt/android/system
    # Mount data
    mount | grep '/dev/mmcblk0p8' || mount /dev/mmcblk0p8 /mnt/android/data

    # --- Provide firmware ---
    [ -L /system ] || ln -s /mnt/android/system /system
    [ -L /data ]   || ln -s /mnt/android/data /data

    # Insmod with wifi support
    modprobe bcmdhd
end script
