# rootfs/etc/init/tf700t.conf
# tf700t

description  "Asus tf700t specific init script"

start on started udev

task

script
    #  mmc read ahead size
    if [ -e /sys/block/mmcblk0/queue/read_ahead_kb ]; then
        echo 2048 > /sys/block/mmcblk0/queue/read_ahead_kb
    fi
    if [ -e /sys/block/mmcblk1/queue/read_ahead_kb ]; then
        echo 2048 > /sys/block/mmcblk1/queue/read_ahead_kb
    fi

    # Disable android lowmemkiller
    if [ -e /sys/module/lowmemorykiller/parameters/minfree ]; then
        echo "0,0,0,0,0,0" > /sys/module/lowmemorykiller/parameters/minfree
    fi

    # --- CPU ---
    # Block lowpower core
    #  it should be enabled on screen off or in strict powersawe mode
    echo 1 > /sys/module/cpu_tegra3/parameters/no_lp

    # Change conservative mode
    echo conservative > /sys/devices/system/cpu/cpu0/cpufreq/scaling_governor
    echo 95 > /sys/devices/system/cpu/cpufreq/conservative/up_threshold
    echo 75 > /sys/devices/system/cpu/cpufreq/conservative/down_threshold
    echo 3 > /sys/devices/system/cpu/cpufreq/conservative/freq_step

    echo interactive > /sys/devices/system/cpu/cpu0/cpufreq/scaling_governor
    echo 2 > /sys/devices/system/cpu/cpufreq/interactive/boost_factor
    echo 1 > /sys/devices/system/cpu/cpufreq/interactive/input_boost
    echo 80 > /sys/devices/system/cpu/cpufreq/interactive/sustain_load

    # CPU hotplugging
    echo 1 > /sys/module/cpu_tegra3/parameters/auto_hotplug

    # lp2 idle state
    echo "Y" > /sys/module/cpuidle/parameters/lp2_in_idle

    # --- GPS ---
    if [ ! -e /sys/block/mmcblk0/queue/read_ahead_kb ]; then
        echo 162 > /sys/class/gpio/export
        #echo 0 > /sys/class/gpio/gpio162/value
        echo out > /sys/class/gpio/gpio162/direction
    fi

    # --- USB ---
    # Wakeup on usb insert
    echo enabled > /sys/bus/usb/devices/1-1/power/wakeup


    # --- Mount android ---
    [ -d /mnt/android/system ] || mkdir -p /mnt/android/system
    [ -d /mnt/android/data ]   || mkdir -p /mnt/android/data
    # Mount system
    mount /dev/mmcblk0p1 /mnt/android/system
    # Mount data
    mount /dev/mmcblk0p8 /mnt/android/data

    # --- Provide firmware ---
    [ -L /system ] || ln -s /mnt/android/system /system
    [ -L /data ]   || ln -s /mnt/android/data /data

    # Insmod with wifi support
    modprobe bcmdhd
end script
